# systemd/waypoint-bus@.service
[Unit]
Description=Waypoint mode bus %i (RFC-0003 / Addendum A)
Documentation=https://github.com/KN4OQW/waypoint/blob/main/docs/rfcs/0003a-loopback-handoff.md
# A DMR bus MULTIPLEXES on DMRGateway (it dials the bus as one [DMR Network N]); a
# YSF/NXDN bus DISPLACES its gateway (binds the mode's loopback). Order after both
# so whatever it rides or replaces is settled first. The mutual exclusion with a
# displaced gateway is enforced by waypointd: the render does not emit the stock
# gateway's target while a bus displaces it, and apply stops it before the bus
# starts and disables it for boot — so no per-instance Conflicts= is needed here (a
# template cannot express one, and the render+apply invariant is stronger).
After=mosquitto.service waypoint-mmdvm.service waypoint-dmrgateway.service
Wants=mosquitto.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
ExecStart=/home/pi-star/waypoint/bin/waypoint-bus -config /home/pi-star/waypoint/etc/waypoint-bus-%i.json -node %H -dmrids /usr/local/etc/DMRIds.dat
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
